# php/_common/prod/fpm-alpine/profile.d/0paths.sh
#!/bin/sh
export PATH="/root/bin:$PATH"

SRC_PATH="/src"
export APP_ROOT="$SRC_PATH"
export APP_LOG="$SRC_PATH/log"
export APP_TMP="$SRC_PATH/tmp"